# deploy/docmanagement.service
[Unit]
Description=Document Management Service
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/docmanagement
Environment="PATH=/var/www/docmanagement/venv/bin"
ExecStart=/var/www/docmanagement/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/docmanagement/access.log
StandardError=append:/var/log/docmanagement/error.log

[Install]
WantedBy=multi-user.target
